# configs/getty@tty1.service.d/autologin.conf
# /etc/systemd/system/getty@tty1.service.d/autologin.conf
# Auto-login __USER__ on the physical console (tty1) so the ops dashboard can
# start unattended after a reboot. setup.sh substitutes __USER__ and installs
# this; the matching ~/.bash_profile hook exec's smart_display.sh on tty1.
#
# Trade-off: anyone at the physical keyboard gets a logged-in shell. Fine for a
# lid-closed box in your space; don't do this on hardware others can reach.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin __USER__ --noclear %I $TERM
