# configs/homelab-health.timer
# /etc/systemd/system/homelab-health.timer
# Run the health check shortly after boot, then every 5 minutes.
[Unit]
Description=Run the Homelab health check on a schedule

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
